# packages/mediacenter/OutFox/system.d/outfox-sextetstream.socket
[Unit]
Description=OutFox SextetStream FIFO
ConditionKernelCommandLine=!installer

[Socket]
ListenFIFO=/tmp/outfox/StepMania-Lights-SextetStream.out

[Install]
WantedBy=outfox.target
